# assets/etc/supervisord/matomo_watch_plugins.conf
[program:matomo_watch_plugins]
command = /usr/local/bin/matomo_watch_plugins
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
